# update/AtualizarControleAcesso.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Modelo de referencia — prefira scripts\agendar_atualizacao.ps1 (gera caminhos corretos).
  Importacao manual (Admin, ajuste caminhos antes):
    schtasks /Create /TN "ControleAcesso-AtualizacaoDiaria" /XML "C:\...\update\AtualizarControleAcesso.xml"
-->
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.4">
  <RegistrationInfo>
    <Description>Atualiza o Controle de Acesso diariamente as 00:00. Log em logs\update.log</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-01T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <RunLevel>HighestAvailable</RunLevel>
      <UserId>S-1-5-18</UserId>
      <LogonType>ServiceAccount</LogonType>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <Enabled>true</Enabled>
    <ExecutionTimeLimit>PT2H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>cmd</Command>
      <Arguments>/c "C:\inetpub\wwwroot\controle-acesso-PAMC\update_app.bat"</Arguments>
      <WorkingDirectory>C:\inetpub\wwwroot\controle-acesso-PAMC</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
